# tools/vagrant/crontab
*   *  *   *   *     /var/www/app/manage.py import
*   *  *   *   *     /var/www/app/manage.py sendmail
