# tools/vagrant/crontab
*   *  *   *   *     /vagrant/envpy /vagrant/app/manage.py import
*   *  *   *   *     /vagrant/envpy /vagrant/app/manage.py sendmail
